# iptoemail.service
[Service]
ExecStart=/home/$USER/iptoemail.sh


[Install]
WantedBy=multi-user.target
