# iptoemail.service
[Service]
ExecStart=/bin/iptoemail.sh


[Install]
WantedBy=multi-user.target
