# lol.service
[Service]
ExecStart=/home/mike/lol.sh


[Install]
WantedBy=multi-user.target
