# lol.service
[Service]
ExecStart=/home/mike/iptoemail.sh


[Install]
WantedBy=multi-user.target
